# experiment/scw/berty.service
[Unit]
Description=Berty Daemon
After=network.target

[Service]
User=root
Group=root
ExecStart=/bin/bash -c '/usr/bin/berty daemon --dht-server --log-level=debug --mdns=false --hop --bind-p2p /ip4/0.0.0.0/tcp/4004,/ip4/0.0.0.0/tcp/80,/ip4/0.0.0.0/tcp/443,/ip4/0.0.0.0/udp/4004/quic --p2p-identity=`/bin/cat /root/pk`'
Restart=always

[Install]
WantedBy=multi-user.target
